# supervisord.conf
[program:nano-llm]
command=.venv/bin/python main.py
directory=%(here)s
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=60
redirect_stderr=true
stdout_logfile=%(here)s/logs/nano-llm.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
environment=PYTHONUNBUFFERED=1
